# cron.d
* * * * * /home/devil/drona/update >/dev/null 2>&1
